# asset-register.service
[Unit]
Description=Employee Asset Management Register (FastAPI)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Employee_Asset_Management_Register
Environment=PATH=/home/ubuntu/Employee_Asset_Management_Register/venv/bin
ExecStart=/home/ubuntu/Employee_Asset_Management_Register/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure

[Install]
WantedBy=multi-user.target
